# configs/systemd/deepseek-v4-flash-ds4.service
[Unit]
Description=DeepSeek V4 Flash (ds4 engine)
Wants=network-online.target
After=network-online.target

[Service]
# The serve script deliberately daemonizes the engine and its OOM watchdog.
# A oneshot unit tracks that managed state without relying on a forking PID guess.
Type=oneshot
RemainAfterExit=yes
User=dsv4
Group=dsv4
Environment=HOME=/home/dsv4
WorkingDirectory=@DSV4_REPO@
RuntimeDirectory=dsv4
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
# This report remains advisory, but preflight failure blocks engine startup;
# the serve script also enforces its authoritative checks on every start.
ExecStartPre=@DSV4_REPO@/scripts/00_preflight.sh --out /tmp/preflight-service.json
ExecStart=@DSV4_REPO@/scripts/20_serve_ds4.sh start
ExecStop=@DSV4_REPO@/scripts/20_serve_ds4.sh stop
TimeoutStartSec=700
TimeoutStopSec=240
# The external dsv4-guard timer supervises health; the OOM watchdog handles kills.
Restart=no

# /home/dsv4 contains the engine build, model artifacts, and logs, so it cannot
# be hidden. It is the only persistent home directory writable by this service.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run/dsv4 /home/dsv4
PrivateTmp=yes
ProtectHome=no
# Defense in depth: this single-host engine needs only loopback and Unix sockets.
IPAddressAllow=127.0.0.0/8 ::1
IPAddressDeny=any
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
