# configs/systemd/deepseek-v4-flash-llamacpp.service
[Unit]
Description=DeepSeek V4 Flash (llama.cpp engine)
Wants=network-online.target
After=network-online.target
# Defense in depth against retry storms after a startup memory/NVRM OOM. This
# StartLimit only catches FAST bursts (>=3 starts within 900 s), and it counts
# ALL starts (including quick successful manual/installer restarts), so it is a
# secondary guard. The PRIMARY circuit breaker is scripts/03_guard.sh, which
# counts CONSECUTIVE unhealthy checks and latches after 3 regardless of how slow
# each blocking restart is (systemd's time window can be evaded by ~600 s
# readiness waits). Clear StartLimit with `systemctl reset-failed`.
StartLimitIntervalSec=900
StartLimitBurst=3

[Service]
# The serve script deliberately daemonizes the engine and its OOM watchdog.
# A oneshot unit tracks that managed state without relying on a forking PID guess.
Type=oneshot
RemainAfterExit=yes
User=dsv4
Group=dsv4
Environment=HOME=/home/dsv4
# Projected-free memory floor for the startup budget gate. 15 GiB keeps a 3 GiB
# buffer above the 12 GiB external watchdog kill line; the 90 GiB model on this
# 119 GiB host lands at ~16-20 GiB free in practice. 16 was too tight once the
# host's non-reclaimable baseline settled after first boot (projected 15.9).
Environment=DSV4_MEM_FLOOR_GIB=15
# Production serves the fused-hyper-connection build (llama.cpp 0dc74e3, PR
# #25585): validated correct (GSM8K dev 97/100, identical to the 32e789fd
# benchmark build; needle retrieval + golden gates pass) and ~14% faster decode
# / ~11-19% faster prefill. Evidence: docs/phase-b-fusion-2026-07-17.md. Remove
# these two lines to revert to the 32e789fd build under the default paths.
Environment=DSV4_SERVER_BINARY=/home/dsv4/llamacpp-project/src/llama.cpp-fusion/build/bin/llama-server
Environment=DSV4_BUILD_MANIFEST=@DSV4_REPO@/configs/build-manifests/llamacpp-fusion.json
# Runtime disk floor for the ExecStartPre preflight. 350 GiB is the setup-phase
# requirement (fetching both weight sets); a loaded engine only writes logs, so
# 100 GiB is an ample runtime floor now that the host disk is shared/fuller.
Environment=DSV4_MIN_ROOT_FREE_GIB=100
# Tuned 2026-07-23 (docs/speed-tuning-2026-07-23.md). Micro-batch 2048: the
# fusion build's measured CUDA compute buffer is 267 MiB at ub=512 (the 35 GiB
# graph-capture peak that motivated ub=256 was the pre-fusion build); ub=2048
# costs ~1 GiB and lifts 19K-prompt prefill from ~208 to ~434 tok/s. The
# launcher charges +2 GiB membudget overhead when ub>512 (DSV4_UBATCH_LARGE).
Environment=DSV4_UBATCH=2048
Environment=DSV4_BATCH=2048
Environment=DSV4_UBATCH_LARGE=1
# GB10 mmap page faults make cold loads pathologically slow (measured here:
# 6-7 min -> ~92 s with --no-mmap). Weights are resident either way (-ngl 999).
Environment=DSV4_NO_MMAP=1
# Draft-model-free n-gram speculative decoding: no memory cost, no effect on
# prose (18.3 tok/s unchanged), +11% on code echoing context, +50% on
# repetitive JSON/tool output (86% draft acceptance). Lossless at temp 0.
Environment=DSV4_SPEC_TYPE=ngram-map-k4v
WorkingDirectory=@DSV4_REPO@
RuntimeDirectory=dsv4
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
# This report remains advisory, but preflight failure blocks engine startup;
# the serve script also enforces its authoritative checks on every start.
ExecStartPre=@DSV4_REPO@/scripts/00_preflight.sh --out /tmp/preflight-service.json
ExecStart=@DSV4_REPO@/scripts/21_serve_llamacpp.sh start
ExecStop=@DSV4_REPO@/scripts/21_serve_llamacpp.sh stop
TimeoutStartSec=700
TimeoutStopSec=240
# The external dsv4-guard timer supervises health; the OOM watchdog handles kills.
Restart=no

# /home/dsv4 contains the engine build and logs, so it cannot be hidden. It is
# the only persistent home directory writable by this service.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run/dsv4 /home/dsv4
PrivateTmp=yes
ProtectHome=no
# Defense in depth: this single-host engine needs only loopback and Unix sockets.
IPAddressAllow=127.0.0.0/8 ::1
IPAddressDeny=any
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
