# configs/systemd/dsv4-authhelper.service
[Unit]
Description=DeepSeek V4 Flash authenticated streaming engine proxy
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=dsv4auth
Group=dsv4auth
EnvironmentFile=/etc/deepseek-v4-flash/env
Environment=HOME=/run/dsv4auth
Environment=PYTHONDONTWRITEBYTECODE=1
RuntimeDirectory=dsv4auth
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
ExecStart=/usr/bin/python3 /usr/local/lib/deepseek-v4-flash/40_auth_helper.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2
TimeoutStopSec=30

# The helper reads its root-owned key and UPSTREAM_PORT from EnvironmentFile,
# accepts Caddy traffic on loopback, and is the only proxy that reaches engine.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run/dsv4auth
PrivateTmp=yes
ProtectHome=no
IPAddressAllow=127.0.0.0/8 ::1
IPAddressDeny=any
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
