# configs/systemd/dsv4-guard.service
[Unit]
Description=DeepSeek V4 Flash engine status guard
After=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/deepseek-v4-flash/env
ExecStart=@DSV4_REPO@/scripts/03_guard.sh

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run/dsv4 /home/dsv4
PrivateTmp=yes
ProtectHome=no
