# configs/systemd/dsv4-guard.service
[Unit]
Description=DeepSeek V4 Flash engine status guard
After=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/deepseek-v4-flash/env
ExecStart=/bin/bash -Eeuo pipefail -c 'case "$$STACK" in ds4) script=/home/bmarti44/spark-deepseek-v4-flash/scripts/20_serve_ds4.sh; unit=deepseek-v4-flash-ds4.service ;; llamacpp) script=/home/bmarti44/spark-deepseek-v4-flash/scripts/21_serve_llamacpp.sh; unit=deepseek-v4-flash-llamacpp.service ;; *) echo "guard: invalid STACK=$$STACK"; exit 1 ;; esac; echo "guard: checking stack=$$STACK"; if /usr/sbin/runuser -u dsv4 -- "$$script" status; then echo "guard: stack=$$STACK healthy"; else echo "guard: stack=$$STACK unhealthy; restarting $$unit"; /usr/bin/systemctl restart "$$unit"; echo "guard: restart completed unit=$$unit"; fi'

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/run/dsv4 /home/dsv4
PrivateTmp=yes
ProtectHome=no
